# project-spec/meta-user/recipes-apps/dac-tcp-server/files/dac-tcp-server.service
[Unit]
Description=TCP server to control the FPGA DAC registers

[Service]
ExecStart=/usr/bin/dac-tcp-server

[Install]
WantedBy=multi-user.target
